# core/etc/lico.supervisor.d/base.ini
[program:lico-core]
command=%(ENV_PYTHON_PREFIX)s/bin/gunicorn --paste %(ENV_LICO_CONFIG_FOLDER)s/lico.paste.ini -w 10 --threads 20 --bind :18080 --log-config %(ENV_LICO_CONFIG_FOLDER)s/lico.logging.d/base.ini --timeout 3600 --keep-alive 65 --capture-output --preload
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=%(ENV_LICO_LOG_FOLDER)s/%(program_name)s.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10

[program:lico-core-worker]
command=%(ENV_PYTHON_PREFIX)s/bin/lico worker --log-level INFO --autoscale 30,10
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=%(ENV_LICO_LOG_FOLDER)s/%(program_name)s.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10

[program:lico-core-beat]
command=%(ENV_PYTHON_PREFIX)s/bin/lico beat --log-config %(ENV_LICO_CONFIG_FOLDER)s/lico.logging.d/base.ini
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=%(ENV_LICO_LOG_FOLDER)s/%(program_name)s.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
